# config/systemd/rocketchat-server.service
[Unit]
Description=Rocket.Chat Server
Requires=network.target rocketchat-mongodb.service
After=network.target rocketchat-mongodb.service

[Service]
Type=simple
Environment=SNAP_COMMON={{ app_data_dir }}
ExecStart={{ app_dir }}/bin/service.server.sh start
ExecStartPost={{ app_dir }}/bin/service.server.sh post-start
Restart=always
RestartSec=10
TimeoutStartSec=600
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=rocketchat
User=rocketchat

[Install]
WantedBy=multi-user.target 